# etc/service/meas_json/run
#!/bin/sh
exec 2>&1
exec /usr/bin/websocketd --staticdir=/home/whytek/meas_web/meas-web --port=8080 /usr/bin/stdbuf -oL /usr/bin/meas_json

